# local/launchd/com.maestro.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Maestro local daemon LaunchAgent.

  Install:
    cp local/launchd/com.maestro.daemon.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.maestro.daemon.plist

  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.maestro.daemon.plist
    rm ~/Library/LaunchAgents/com.maestro.daemon.plist

  Secrets live in ~/Library/Application Support/maestro/env. Never commit.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.maestro.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/nathancolestock/maestro/local/launchd/run-daemon.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/nathancolestock/maestro/local</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>NetworkState</key>
    <true/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>/Users/nathancolestock/Library/Logs/suite/maestro-daemon.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nathancolestock/Library/Logs/suite/maestro-daemon.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>
</dict>
</plist>
